# pi-setup/systemd/antikythera-gpio-button.service
[Unit]
Description=Antikythera GPIO Button Handler
Documentation=https://github.com/your-repo/antikythera-engine-2
# Run very early in boot, before Comitup
DefaultDependencies=no
Before=comitup.service network-pre.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/antikythera-gpio-button.py
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
