# pi-setup/systemd/antikythera.service
[Unit]
Description=Antikythera Astronomical Engine
Documentation=https://github.com/your-repo/antikythera-engine-2
# Wait for network and WiFi window to close
After=network-online.target antikythera-wifi-window.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/antikythera-engine-2
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/antikythera-engine-2/config

[Install]
WantedBy=multi-user.target
